# infra/systemd/platform-paas-backend.service
[Unit]
Description=Platform PaaS Backend
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=platform
Group=platform
WorkingDirectory=/opt/platform_paas/backend
EnvironmentFile=/opt/platform_paas/.env
ExecStart=/opt/platform_paas/platform_paas_venv/bin/python -m uvicorn app.main:app --host 127.0.0.1 --port 8000 --workers 5
Restart=always
RestartSec=5
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
