# infra/systemd/platform-paas-provisioning-worker@.service
[Unit]
Description=Run provisioning worker once for profile %i
After=postgresql.service
Wants=postgresql.service

[Service]
Type=oneshot
User=platform
Group=platform
WorkingDirectory=/opt/platform_paas
EnvironmentFile=/opt/platform_paas/.env
ExecStart=/opt/platform_paas/deploy/run_provisioning_worker_profile.sh %i
